# terabox-bot.service
[Unit]
Description=Terabox Automation Bot
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/terabox-automation-bot
Environment="PATH=/opt/terabox-automation-bot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/terabox-automation-bot/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=terabox-bot

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
